# mpss3/mpss-modules/udev-mic.rules
# do not edit this file, it will be overwritten on update
# initramfs:default

# MIC SCIF
KERNEL=="scif", ACTION=="add", NAME="mic/%k",MODE="0666", RUN+="/bin/chmod og+x /dev/mic"
KERNEL=="ctrl", ACTION=="add", NAME="mic/%k", MODE="0666"

# Bring up network interfaces manually on rhel7 after module reload
KERNEL=="mic*", SUBSYSTEM=="net", RUN+="/bin/sh -c '/bin/grep 7. /etc/redhat-release && /sbin/ifup %k'"
